# templates/etc/profile.d/sbt.sh.j2
export PATH=$PATH:{{ sbt_usr_dir }}/bin/